# home/.crontab
# Custom cron

#
# This ensures that a custom mimeapps.list file is always used.
# Sometimes there are apps that replace this file with whatever they want.
#
*/5 * * * * ln -fs $HOME/dotfiles/home/.config/mimeapps.list $HOME/.config/.

#
# Backups
#
# Run daily backup at 00:00 AM
0 0 * * * /home/cmiranda/Documents/Work/Workspace/dotfiles/home/bin/backup backup --quiet
#
# Keep only the last 14 daily backups (run at 4:00 AM)
0 4 * * * /home/cmiranda/Documents/Work/Workspace/dotfiles/home/bin/backup prune --retention 14 --quiet

#
# Creates a dump of the local ESH Oracle database
#
@monthly $HOME/bin/esh-db-backup
