# home/.crontab
# Custom cron

*/5 * * * * bash /home/$USER/bin/set-random-wallpaper
