# home/.crontab
# Custom cron

#
# This ensures that a custom mimeapps.list file is always used.
# Sometimes there are apps that replace this file with whatever they want.
#
*/5 * * * * ln -fs /home/$USER/dotfiles/home/.config/mimeapps.list /home/$USER/.config/.

@weekly bash /home/$USER/bin/backup
