# box/linux-mint/cron/root.cron
@hourly /home/cmiranda/bin/backuper --backup linux
